# home/Library/LaunchAgents/net.aaroncurry.motherbox.nightly-maintenance.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>net.aaroncurry.motherbox.nightly-maintenance</string>

	<key>ProgramArguments</key>
	<array>
		<string>{{ .chezmoi.sourceDir | dir }}/scripts/utils/maintenance.sh</string>
		<string>{{ .machine }}</string>
	</array>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:{{ .chezmoi.homeDir }}/.local/bin</string>
	</dict>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>2</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<key>StandardOutPath</key>
	<string>{{ .chezmoi.homeDir }}/.config/motherbox/logs/nightly-maintenance.log</string>
	<key>StandardErrorPath</key>
	<string>{{ .chezmoi.homeDir }}/.config/motherbox/logs/nightly-maintenance.log</string>
</dict>
</plist>
